# xsession/lock.service
[Unit]
Description=Lock the screen on resume from suspend
Before=sleep.target

[Service]
User=user_placeholder
Environment=DISPLAY=:0
ExecStartPre=/usr/bin/xset dpms force suspend
ExecStart=/usr/local/bin/slock

[Install]
WantedBy=sleep.target hibernate.target suspend.target suspend-then-hibernate.target hybrid-sleep.target
